# monitoring/roles/monitoring/files/node_exporter.service
# Managed by Ansible - do not modify

[Unit]
Description=node_exporter
After=network.target

[Service]
ExecStart=/usr/local/bin/node_exporter
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
